# scripts/deploy/callmanagement-worker.service
[Unit]
Description=Call Management LiveKit Agent Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/callmanagement
EnvironmentFile=/opt/callmanagement/.env
ExecStart=/home/ubuntu/.local/bin/uv run -m call_management.server start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target